# cron/crontab
* * * * * root cd /app && /usr/local/bin/python manage.py deactivate_urls >> /var/log/cron.log 2>&1